# backend/convince-backend.service
[Unit]
Description=Convince AI Backend
After=network.target

[Service]
Type=exec
User=aloo
Group=aloo
WorkingDirectory=/home/aloo/convince-ai/backend
Environment=PATH=/home/aloo/convince-ai/venv/bin
Environment=PYTHONPATH=/home/aloo/convince-ai/backend
ExecStart=/home/aloo/convince-ai/venv/bin/gunicorn -c gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=convince-backend

# Security settings
NoNewPrivileges=true
DynamicUser=false
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/aloo/convince-ai/backend
PrivateTmp=true

[Install]
WantedBy=multi-user.target
